# .jdocker.cron
## jdocker cron
jdocksh=DIR/jdocker.sh

# Exemple
#0 4 * * 0   root $jdocksh bk exemple f >/dev/null 2>&1
#0 4 * * 1-6 root $jdocksh bk exemple i >/dev/null 2>&1

# Minecraft
0 7 * * * root $jdocksh up mcserver >/dev/null 2>&1
