# .jdocker.cron.template
# jdocker cron
jdocksh=$dir/jdocker.sh

# exemple1
0 0 * * *  $user $jdocksh bk exemple1 >/dev/null 2>&1

# exemple2
0 1 * * *  $user $jdocksh bk exemple2 >/dev/null 2>&1
