# jdocker.cron
# jdocker cron
jdocksh=DIR/jdocker.sh

# nginx
56 20 * * *  root $jdocksh bk nginx >/dev/null 2>&1
