# jdocker.cron
# jdocker cron
user=USER
jdocksh=DIR/jdocker.sh

# Exemple
0 4 * * *  $user $jdocksh bk exemple >/dev/null 2>&1
