# jdocker.cron
# jdocker cron
jdocksh=SCR

# mcserver
0 7 * * *  USER $jdocksh bk mcserver >/dev/null 2>&1
