# jdocker.cron
# jdocker cron
jdocksh=SCR

# exemple
0 0 * * *  root $jdocksh bk exemple >/dev/null 2>&1
