# jdocker.cron
# jdocker cron
jdocksh=DIR/jdocker.sh

# exemple
0 0 * * *  root $jdocksh bk exemple >/dev/null 2>&1
