# jdocker.cron
# jdocker cron
jdocksh=DIR/jdocker.sh

# example
0 0 * * *  USER $jdocksh bk exemple >/dev/null 2>&1
