# jdocker.cron
# jdocker cron
jdocksh=SCR

# exemple1
0 0 * * *  USER $jdocksh bk exemple1

# exemple2
0 1 * * *  USER $jdocksh bk exemple2
