# jdocker.cron
# jdocker cron
jdocksh=DIR/jdocker.sh

# emulatorjs
18 23 * * *  USER $jdocksh bk emulatorjs >/dev/null 2>&1
